# apps/mongodb/mongodb.service
[Unit]
Description=MongoDB
After=docker.service
Requires=docker.service

[Service]
ExecStart=/bin/bash -c '/usr/bin/docker start -a mongodb || /usr/bin/docker run --name mongodb wekeypedia/mongodb'
ExecStop=/usr/bin/docker stop mongodb 

[Install]
WantedBy=multi-user.target